# simplewebapp.service
[Service]
Type=simple
Environment="PORT=5000"